# server/ztna-server.service
[Unit]
Description=ZTNA DTLS Server (Control Plane + Data Plane)
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/ztna

# 設定ファイルの事前検証
ExecStartPre=/opt/ztna/venv/bin/python3 -c "import yaml; yaml.safe_load(open('/opt/ztna/config.yaml'))"

# コントロールプレーン(Flask)とデータプレーン(DTLSゲートウェイ)を同時起動
# ── foreman的なランチャースクリプトを使う方式 ──
ExecStart=/bin/bash -c '\
  /opt/ztna/venv/bin/python3 /opt/ztna/app.py &\
  APP_PID=$!;\
  /opt/ztna/venv/bin/python3 /opt/ztna/dtls_gateway.py &\
  GW_PID=$!;\
  echo "started app(PID=$APP_PID) gateway(PID=$GW_PID)";\
  wait'

ExecStop=/bin/bash -c '\
  pkill -TERM -f "python3 /opt/ztna/app.py" || true;\
  pkill -TERM -f "python3 /opt/ztna/dtls_gateway.py" || true'

KillMode=process
Restart=on-failure
RestartSec=5s

# TUNデバイス操作のためにケイパビリティを付与
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ztna-server

[Install]
WantedBy=multi-user.target
